# recycle.supervisor.conf
[program:recycle]
command=/app/service/recycle/recycle
directory=/app/service/recycle
autostart=true
redirect_stderr=true
stdout_logfile=/app/service/recycle/std.out